# output/com.tony.xz-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tony.xz-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tony/dev/personal/xz/update.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tony/dev/personal/xz</string>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/tony/dev/personal/xz/output/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tony/dev/personal/xz/output/launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/tony/.local/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
